# packer/drone-server-configs/drone_reload.service
[Unit]
Description=Refreshes drone images and updates containers

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload-or-restart drone.service
